# systemd/web-server.service
[Unit]
Description=ACKMUD web server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/user/web
ExecStart=/usr/bin/python3 /home/user/web/web_who_server.py
Environment=HOME=/home/user
Environment=ACK_WEB_PORT=8080
Environment=ACKTNG_GAME_URL=http://192.168.1.103:8080
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
